# supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn -w 4 -k uvicorn.workers.UvicornWorker main:app
directory=/app
user=root

[program:bot]
command=python bot.py
directory=/app
user=root
